# packages/printf/service/printf/finish
#! /bin/sh

/opt/00admin/sbin/fire-ip del
