# sh4pwn/skeleton/service/eth0/run
#! /bin/sh -e

ip link set dev eth0 up

exec inotifyd true $(pwd)

