# mcp/service/mathopd/run
#! /bin/sh

exec mathopd -n -f ./mathopd.conf >> /var/log/mathopd.log 2>> /var/log/mathopd.error.log
