# packages/tanks/service/tanksd/run
#! /bin/sh

exec 2>&1

mkdir -p /var/www/tanks
ln -s /opt/tanks/html/* /var/www/tanks/
ln -sf summary.html /var/www/tanks/index.html

install -o ctf -d /var/lib/ctf/tanks/players
exec ./tanksd
