# packages/tanks/service/tanksd/run
#! /bin/sh -e

exec 2>&1

mkdir -p $CTF_BASE/state/tanks/players
mkdir -p $CTF_BASE/state/tanks/www
chown ctf $CTF_BASE/state/tanks
chown ctf $CTF_BASE/state/tanks/players

ln -s $CTF_BASE/state/tanks/www $CTF_BASE/www/tanks
ln -s $CTF_BASE/packages/tanks/html/* $CTF_BASE/state/tanks/www/ || true
ln -s summary.html $CTF_BASE/tanks/www/index.html || true

PATH=/bin:$CTF_BASE/packages/tanks/bin; export PATH
exec ./tanksd
