# packages/logger/service/logger/run
#! /bin/sh

exec 2>&1
exec tcpsvd 0 4104 /opt/logger/bin/logger
