# packages/logger/service/logger/run
#! /bin/sh -e

exec 2>&1
IP=$(cat ip.txt)
ip addr add $IP label eth0:logger dev eth0
exec tcpsvd ${IP#/*} 1958 /opt/logger/bin/logger
