# packages/rlyeh/service/rlyeh/run
#! /bin/sh -e

exec 2>&1

IP=$(/opt/00admin/sbin/fire-ip add)

dir=/var/lib/ctf/rlyeh
install -o nobody -d $dir

PATH=/opt/ctfbase/bin:/opt/rlyeh/bin:/bin:/sbin; export PATH
exec tcpsvd ${IP%/*} 1928 ./rlyeh-ctf
