# packages/rlyeh/service/rlyeh/run
#! /bin/sh -e

exec 2>&1
read IP < ip.txt
ip addr add $IP label eth0:rlyeh dev eth0

dir=/var/lib/ctf/rlyeh
install -o nobody -d $dir

PATH=/opt/ctfbase/bin:/opt/rlyeh/bin:/bin:/sbin; export PATH
exec tcpsvd ${IP%/*} 1928 ./rlyeh-ctf
