# packages/rlyeh/service/rlyeh/run
#! /bin/sh -e

exec 2>&1
read IP < ip.txt
ip addr add $IP label eth0:rlyeh dev eth0

dir=/var/lib/ctf/rlyeh
install -o nobody -d $dir

exec setuidgid nobody tcpsvd ${IP#/*} 1928 ./rlyeh-ctf
