# packages/pwnables/service/pwnables/finish
#! /bin/sh

IP=$(cat ip.txt)
ip addr del $IP dev eth0
