# packages/playfair/service/playfair/run
#! /bin/sh -e

exec 2>&1

IP=$(/opt/00admin/bin/fire-ip add)
ip addr add $IP dev eth0 || true

exec tcpsvd -u nobody ${IP%/*} 1013 /opt/playfair/bin/playfair

