# packages/playfair/service/playfair/run
#! /bin/sh -e

exec 2>&1

IP=$(dbip -a)
ip addr add $IP dev eth0 || true

exec tcpsvd -u nobody ${IP%/*} 1013 $CTF_BASE/packages/playfair/bin/playfair

