# pwnables/skel/var/service/telnetd/run
#! /bin/sh

state=$(cat /sys/class/net/eth1/operstate)
if [ "$state" != "up" ]; then
    exit 0
fi
exec tcpsvd 192.168.0.1 23 telnetd -i -l /bin/sh
