# packages/mcp/service/ircd/run
#! /bin/sh

exec 2>&1
exec ngircd --config ./ngircd.conf --nodaemon
