# packages/p2client/service/eth0/run
#! /bin/sh

exec 2>&1

exec udhcpc -f eth0
