# packages/p2client/service/eth0/run
#! /bin/sh

exec 2>&1

if grep 'debug' /proc/cmdline; then
    ip link set dev eth0 up
    ip addr add dev eth0 10.0.0.110/24
    exec sleep 8d
else
    exec udhcpc -f eth0
fi
