# mcp/service/tokend/run
#! /bin/sh

exec 2>&1
exec tcpsvd -u ctf 0 1 /opt/mcp/bin/in.tokend
