# packages/router/service/sshd/run
#! /bin/sh

exec 2>&1

cat passwd > /etc/passwd

exec dropbear -r ./rsa.key -E -F
