# packages/router/service/sshd/run
#! /bin/sh

exec 2>&1

echo 'root:$1$xAJ7KwiU$BeKJjYGs9r/hY9Ag4qv4I1:0:0:root:/tmp:/bin/sh' > /etc/passwd

exec dropbear -r ./rsa.key -E -F
