# packages/router/service/radvd/run
#! /bin/sh

exec 2>&1

exec /opt/router/bin/radvd -C radvd.conf -d 1 -m stderr -p radvd.pid
