# sh4pwn/skeleton/service/tty1/run
#! /bin/sh

dev=/dev/$(basename $(pwd))
test -c $dev || exit

exec getty 38400 $dev
