# sh4pwn/skeleton/service/telnetd/run
#! /bin/sh

exec tcpsvd :: 23 telnetd -f issue -l /bin/sh
