# packages/pwnables/service/eth0.pwn/run
#! /bin/sh

if [ -d /opt/mcp ]; then
    sv d .
    exit
fi

hostname pwnables
ifconfig eth0 10.0.0.10 netmask 255.255.0.0
exec inotifyd true $(pwd):x
