# packages/ctfbase/service/sshd/finish
#! /bin/sh

iptables -D INPUT -s 10.0.0.0/16 --proto tcp --dport 55 -j ACCEPT
iptables -D INPUT --proto tcp --dport 55 -j REJECT
