# packages/router/service/router/finish
#! /bin/sh

for i in $(seq 24); do
  ip link del link dev eth0.$i
done

ip link set dev eth0 down

iptables -F INPUT
