# packages/mcp/service/sshd/run
#! /bin/sh

exec 2>&1
exec dropbear -r ./rsa.key -E -F
