# packages/p2/service/httpd/finish
#! /bin/sh

/opt/00admin/bin/fire-ip del
