# packages/p2/service/httpd/run
#! /bin/sh -e

exec 2>&1

IP=$(dbip -p ../p2console/ip.txt)

ip route add default via 10.0.0.1

# Link in puzzles and web pages
install -d /var/www
for d in /opt/*; do
    w=/var/www/$(basename $d)
    if [ -d $d/puzzles ] && ! [ -d $w ]; then
        ln -sf $d/puzzles $w
    fi
    if [ -d $d/www ]; then
        ln -sf $d/www/* /var/www/
    fi
done

cd /var/www
exec tcpsvd -u nobody ${IP%/*} 80 /opt/p2/bin/eris -c.
