# packages/p2client/service/p2client/run
#! /bin/sh

exec 2>&1

# Turn off conflicting getty
sv d tty1

if ! sv s tty1 | grep down; then
    exit
fi

# Cool font
setfont -C /dev/tty1 $CTF_BASE/packages/p2client/lite-16.fnt


if ! grep -q debug /proc/cmdline; then
    # Don't let C-A-F2 switch away
    sv d tty2 tty3 tty4
    deallocvt

    # Disable magic sysrq key
    echo 0 > /proc/sys/kernel/sysrq

    # Turn off screen blanking
    setterm -blank 0 -powersave off -powerdown 0

    # Make this like a VT52; including disabling ctrl-alt-del
    loadkmap < $CTF_BASE/packages/p2client/dumbterm.kmap

    # Disable console logging
    echo 1 4 1 4 > /proc/sys/kernel/printk
fi

# Run the console!
exec getty -n -l ./p2connect 38400 tty1

