# packages/p2client/service/p2client/run
#! /bin/sh

exec 2>&1

# Turn off default gettys
sv d tty1 tty2 tty3 tty4
if ! sv s tty1 | grep down; then
    exit
fi

# Don't let C-A-F2 switch away
deallocvt

# Disable magic sysrq key
echo 0 > /proc/sys/kernel/sysrq

# Make this like a VT52; including disabling ctrl-alt-del
loadkmap < /opt/p2client/dumbterm.kmap

# Run the console!
exec getty -n -l ./p2connect 38400 tty1

