# packages/octopus/service/octopus/run
#! /bin/sh -e

exec 2>&1

IP=$(cat ip.txt)
ip addr add $IP dev eth0 || true

exec /opt/octopus/bin/octopus ${IP%/*} < /opt/octopus/tokens.txt
