# packages/octopus/service/octopus/run
#! /bin/sh -e

exec 2>&1

IP=$(dbip -a)

grep -q ipv4 /proc/cmdline && IP=::ffff:$IP

exec /opt/octopus/bin/octopus ${IP%/*} < /opt/octopus/tokens.txt
