# packages/wopr/service/wopr/run
#! /bin/sh

# Work around busybox bug
chmod a+r /dev/urandom

mkdir -p /var/tmp/wopr
chown ctf /var/tmp/wopr
