# packages/cowbull/service/cowbull/finish
#! /bin/sh

read IP < ip.txt
ip addr del $IP dev eth0
