# packages/cowbull/service/cowbull/run
#! /bin/sh -e

exec 2>&1

IP=$(/opt/00admin/bin/fire-ip add)

exec setuidgid nobody /opt/cowbull/bin/cowd < /opt/cowbull/tokens.txt
