# packages/cowbull/service/cowbull/run
#! /bin/sh -e

exec 2>&1
read IP < ip.txt
ip addr add $IP dev eth0 || true

exec setuidgid nobody /opt/cowbull/bin/cowd < /opt/cowbull/tokens.txt
