# packages/cowbull/service/cowbull/run
#! /bin/sh -e

exec 2>&1

IP=$(dbip -a)

exec setuidgid nobody /opt/cowbull/bin/cowd < /opt/cowbull/tokens.txt
