# pwnables/service/pwnables/run
#! /bin/sh -e

exec tcpsvd -C 5:"Let's not be greedy" 0 23 /sbin/telnetd -l ./pwnie
