# packages/octopus/service/octopus-redirect/run
#! /bin/sh -e

IP=$(cat ../octopus/ip.txt)
sv s octopus >/dev/null || exit 1
exec tcpsvd ${IP#/*} 8888 ./octopus-redirect
