# packages/octopus/service/octopus-redirect/run
#! /bin/sh -e

exec 2>&1

IP=$(dbip -p ../octopus/ip.txt)

sv s octopus >/dev/null || exit 1
exec tcpsvd ${IP%/*} 8888 ./octopus-redirect
