# packages/mcp/service/httpd/run
#! /bin/sh

exec 2>&1

install -d /var/www

# Link in puzzles and web pages
for d in /opt/*; do
    if [ -d $d/puzzles ]; then
        ln -sf $d/puzzles /var/www/$(basename $d)
    fi
    if [ -d $d/www ]; then
        ln -sf $d/www/* /var/www/
    fi
done

exec mathopd -n -f ./mathopd.conf

