# packages/mcp/service/httpd/run
#! /bin/sh -e

exec 2>&1

IP=$(/opt/00admin/sbin/fire-ip add)

hostname mcp

# Link in puzzles and web pages
install -d /var/www
for d in /opt/*; do
    w=/var/www/$(basename $d)
    if [ -d $d/puzzles ] && ! [ -d $w ]; then
        ln -sf $d/puzzles $w
    fi
    if [ -d $d/www ]; then
        ln -sf $d/www/* /var/www/
    fi
done

ln -sf /var/www default
exec tcpsvd -u ctf ${IP%/*} 80 /opt/mcp/bin/eris -c
