# packages/mcp/service/httpd/run
#! /bin/sh

exec 2>&1

# Bring up address
IP=$(cat ip.txt)
ip addr add $IP dev eth0
ip monitor | grep -q $IP

# Link in puzzles and web pages
install -d /var/www
for d in /opt/*; do
    if [ -d $d/puzzles ]; then
        ln -sf $d/puzzles /var/www/$(basename $d)
    fi
    if [ -d $d/www ]; then
        ln -sf $d/www/* /var/www/
    fi
done


mkdir -p sites
cd sites
ln -sf /var/www default
exec tcpsvd -u ctf $IP 80 /opt/mcp/bin/fnord-cgi
