# mcp/service/puzzled/run
#! /bin/sh -e

exec 2>&1

DB=/var/lib/ctf/puzzles.db

mkdir -p /var/www
touch $DB
chown ctf $DB

# Make an initial listing
./puzzled

exec inotifyd ./puzzled ${DB}:w
