# packages/ircd/service/ircd/run
#! /bin/sh

exec 2>&1

IP=$(cat ip.txt)
ip addr add $IP label eth0:ircd dev eth0
ip monitor | grep -q $IP

adduser -S -H -u 65534 nobody
adduser -S -H irc

exec setuidgid irc /opt/ircd/bin/ngircd --config ./ngircd.conf --nodaemon
