# packages/ircd/service/ircd/run
#! /bin/sh

exec 2>&1
exec /opt/ngircd/bin/ngircd --config ./ngircd.conf --nodaemon
