# packages/00admin/service/sshd/run
#! /bin/sh

exec 2>&1

echo 'root:$1$i/EeZyl6$5C6Z8s0ftEdTKFLUli1wP1' | chpasswd --encrypted

# Bring up the NIC; this will get us at least a link-local address, and
# hopefully a global address with stateless autoconfiguration.
ip link set eth0 up

exec dropbear -r ./rsa.key -E -F
