# packages/inferno/service/inferno/run
#! /bin/sh

exec 2>&1

mkdir -p /mnt/inferno
mount -o ro LABEL=inferno /mnt/inferno
ln -s /mnt/inferno /var/www

ip addr add eth0 10.0.0.2/24

exec inotifyd true $PWD
