# packages/multicaster/service/multicaster/run
#! /bin/sh -e

exec 2>&1

IP=$(dbip -a)

exec setuidgid ctf /opt/multicaster/bin/multicaster ff15::62c 1580 </opt/multicaster/tokens.txt
