# packages/00common/service/sshd/run
#! /bin/sh

exec 2>&1

password='grape guts'
for fn in /opt/*/password; do
    read password < $fn && break
done

printf 'root:%s' "$password" | chpasswd --md5

# Bring up the NIC; this will get us at least a link-local address, and
# hopefully a global address with stateless autoconfiguration.
ip link set eth0 up

exec dropbear -r ./rsa.key -E -F
